# config/cron.d/dxf_export
0 1 * * * root /app/scripts/command.sh >> /var/log/dxf_cron.log 2>&1